# cosmos-bot.service
[Unit]
Description=Cosmos Monitor Telegram Bot
After=network.target

[Service]
Type=simple
User=your_linux_user
WorkingDirectory=/home/your_linux_user/cosmos-monitor-bot
Environment=TELEGRAM_BOT_TOKEN=YOUR_TOKEN_HERE
Environment=TELEGRAM_CHAT_ID=YOUR_CHAT_ID_HERE
ExecStart=/usr/bin/node bot.js
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
